# root/etc/services.d/flask-websockets/run
#!/usr/bin/with-contenv bash

echo "Starting gunicorn app Server."

cd /app/

exec \
	s6-setuidgid abc gunicorn --worker-class eventlet -w 10 main:app
