# root/etc/services.d/flask-websockets/run
#!/usr/bin/with-contenv bash

echo "Starting flask webserver."

cd /app/

exec \
	/usr/bin/python3 main.py
