# root/etc/services.d/flask-websockets/run
#!/usr/bin/with-contenv bash

echo "Starting gunicorn app Server."

cd /app/

exec \
  python3 mqtt.py
