# root/etc/services.d/redis-server/run
#!/usr/bin/with-contenv bash

echo "Starting redis-server."

exec \
	s6-setuidgid redis /usr/bin/redis-server /etc/redis/redis.conf
